# deploy/ubuntu/unpack-module.service
[Unit]
Description=Unpack Module web application
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=unpack-module
Group=unpack-module
WorkingDirectory=/opt/unpack-module
EnvironmentFile=/opt/unpack-module/.env
ExecStart=/opt/unpack-module/.venv/bin/gunicorn --chdir web --workers 1 --threads 4 --bind 127.0.0.1:5050 --access-logfile - --error-logfile - wsgi:app
Restart=on-failure
RestartSec=5
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target

